# installer/Payload/Library/LaunchAgents/com.yabai.tomodachi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yabai.tomodachi</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/YabaiRestarter</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>Comment</key>
    <string>Manual start only - use 'launchctl load ~/Library/LaunchAgents/com.yabai.tomodachi.plist' to enable auto-start</string>
    <key>LSUIElement</key>
    <true/>
</dict>
</plist>